# scripts/mcp-knowledge-base.service
[Unit]
Description=MCP Knowledge Base Server
Documentation=https://github.com/example/mcp-knowledge-base
After=network.target

[Service]
Type=simple
User=mcp-kb
Group=mcp-kb
WorkingDirectory=/opt/mcp-knowledge-base

# Environment variables for configuration
Environment="LANCE_DB_PATH=/var/lib/mcp-kb/knowledge_base.db"
Environment="LANCE_DB_TABLE_PREFIX=kb_"
Environment="EMBEDDING_MODEL=sentence-transformers/all-MiniLM-L6-v2"
Environment="EMBEDDING_BATCH_SIZE=32"
Environment="EMBEDDING_CACHE_DIR=/var/lib/mcp-kb/models"
Environment="MAX_FILE_SIZE_MB=10"
Environment="CHUNK_SIZE=1000"
Environment="CHUNK_OVERLAP=200"
Environment="SUPPORTED_EXTENSIONS=.py,.txt"
Environment="MCP_SERVER_NAME=knowledge-base"
Environment="LOG_LEVEL=INFO"

# Alternatively, use a configuration file
# ExecStart=/usr/local/bin/mcp-knowledge-base --config /etc/mcp-kb/config.yaml --log-dir /var/log/mcp-kb

# Start command
ExecStart=/usr/local/bin/mcp-knowledge-base --log-dir /var/log/mcp-kb

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/mcp-kb /var/log/mcp-kb

# Resource limits
LimitNOFILE=65536
MemoryLimit=2G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp-knowledge-base

[Install]
WantedBy=multi-user.target
